# specsavers-checker.timer
[Unit]
Description=Checks for free SpecSavers appointments every 15 minutes (Mon-Sat, 9 AM - 5:30 PM)

[Timer]
# Every 15 minutes from 10:00 AM to 11:45 PM
OnCalendar=10..23:00,15,30,45
Persistent=true

[Install]
WantedBy=timers.target
